# config/etc/services.d/mysql/run
#!/bin/sh
if [ ! -d "/run/mysqld" ]; then
  mkdir -p /run/mysqld
fi
/usr/bin/mysqld --user=root --console
